# docs/example/thoth.socket
[Unit]
Description=thoth socket

[Socket]
ListenStream=/run/thoth.sock

[Install]
WantedBy=sockets.target